# services/zpanel-app-notify@.service
[Unit]
Description=ZPanel App Status Notifier for %i

[Service]
Type=oneshot
ExecStart=/zpanel/bin/zapp-notify.sh %i